# sshail/conf/etc/cron.d/sshail
# Call the sshail purge every hour

19 * * * *   sshail  [ -x /usr/local/bin/sshail ] && /usr/local/bin/sshail --purge
